# bin/hs_serv/supervisor/conf.d/cs360.conf
[program:cs3601] 
directory=/home/zhou/sites/cs360 
command=python cs360.py -port=5010
autostart=true 
autorestart=true 

[program:cs3602] 
directory=/home/zhou/sites/cs360
command=python cs360.py -port=5011
autostart=true 
autorestart=true 

[program:cs3603] 
directory=/home/zhou/sites/cs360
command=python cs360.py -port=5012
autostart=true 
autorestart=true 

[program:cs3604] 
directory=/home/zhou/sites/cs360
command=python cs360.py -port=5013
autostart=true 
autorestart=true 

[program:cs3605] 
directory=/home/zhou/sites/cs360
command=python cs360.py -port=5014
autostart=true 
autorestart=true 

[program:cs3606] 
directory=/home/zhou/sites/cs360
command=python cs360.py -port=5015
autostart=true 
autorestart=true 

[program:cs3607] 
directory=/home/zhou/sites/cs360
command=python cs360.py -port=5016
autostart=true 
autorestart=true 

[program:cs3608] 
directory=/home/zhou/sites/cs360
command=python cs360.py -port=5017
autostart=true 
autorestart=true 
