# bin/hs_serv/supervisor/conf.a/esl_go_9090.conf
[program:esl_go_9090] 
directory=/home/zhou
command= /home/zhou/bin/new9090 
autostart=true 
autorestart=true 

