# bin/hs_serv/supervisor/conf.a/esl_5.conf
[program:esl5] 
directory=/home/zhou/bin/eventsocket 
command= twistd --pidfile=server5.pid -ny server5.py
autostart=true 
autorestart=true 

